# build/dotnetcore/hellomvc.conf
[program:hellomvc]
command=/usr/bin/dotnet HostApi.dll --server.urls=http://*:5000
directory=/var/aspnetcore/HelloMVC/publish
autostart=true
autorestart=true
stderr_logfile=/var/log/hellomvc.err.log
stdout_logfile=/var/log/hellomvc.out.log
environment=HOME=/var/www/,ASPNETCORE_ENVIRONMENT=Production,COREHOST_TRACE=1
user=www-data
stopsignal=INT
stopasgroup=true
killasgroup=true
